# bckp/crontab
SHELL=/bin/sh
MAILTO=aleister

# Auto suspend laptop
* * * * *	aleister	~/.local/scripts/bat
